# contrib/init/muskcoind.conf
description "Muskcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env MUSKCOIND_BIN="/usr/bin/muskcoind"
env MUSKCOIND_USER="muskcoin"
env MUSKCOIND_GROUP="muskcoin"
env MUSKCOIND_PIDDIR="/var/run/muskcoind"
# upstart can't handle variables constructed with other variables
env MUSKCOIND_PIDFILE="/var/run/muskcoind/muskcoind.pid"
env MUSKCOIND_CONFIGFILE="/etc/muskcoin/muskcoin.conf"
env MUSKCOIND_DATADIR="/var/lib/muskcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # muskcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$MUSKCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run muskcoind."
        echo "The setting must appear in $MUSKCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Muskcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$MUSKCOIND_PIDDIR"
    chmod 0755 "$MUSKCOIND_PIDDIR"
    chown $MUSKCOIND_USER:$MUSKCOIND_GROUP "$MUSKCOIND_PIDDIR"
    chown $MUSKCOIND_USER:$MUSKCOIND_GROUP "$MUSKCOIND_CONFIGFILE"
    chmod 0660 "$MUSKCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$MUSKCOIND_PIDFILE" \
    --chuid $MUSKCOIND_USER:$MUSKCOIND_GROUP \
    --exec "$MUSKCOIND_BIN" \
    -- \
    -pid="$MUSKCOIND_PIDFILE" \
    -conf="$MUSKCOIND_CONFIGFILE" \
    -datadir="$MUSKCOIND_DATADIR" \
    -disablewallet \
    -daemon

